# ops/systemd/imperium-core-backup.timer
[Unit]
Description=Run Imperium Core PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 03:15:00 UTC
Persistent=true
RandomizedDelaySec=15m
Unit=imperium-core-backup.service

[Install]
WantedBy=timers.target
